# packages-initramfs/rc.local/etc/rc.local
#!/bin/sh -e

. /etc/init.d/oops

cat /proc/device-tree/compatible | cut -d, -f2 | \
while read machine; do run-parts --arg start --exit-on-error /usr/share/rc.local/$machine.d && break; done

export TMPDIR=/tmp

for dir in /usr/local/sbin /usr/local/bin; do
	[ -d $dir ] && PATH=$dir:$PATH
done
